# review_system/cron/crontab
0 8 * * * cd /app/app && /usr/local/bin/python main.py >> /var/log/cron.log 2>&1
0 18 * * * cd /app/app && /usr/local/bin/python main.py >> /var/log/cron.log 2>&1
